# dist/thinkhdd.service
[Unit]
Description=ThinkHDD i LED control

[Service]
Type=simple
ExecStart=/usr/bin/thinkhdd
